# src/biostimer-agent-ms-cleaner.timer.in
[Unit]
Description=Timer to call biostimer-agent-ms-cleaner.service
PartOf=bios.target

[Timer]
### Run every midnight
OnCalendar=*-*-* 23:37:00
# Do not record last-execution times
Persistent=false

[Install]
WantedBy=bios.target
