# src/bios-agent-ms.service.in
[Unit]
Description=BIOS metric persistance
After=malamute.service
PartOf=bios.target

[Service]
Type=simple
Environment="prefix=@prefix@"
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/bios/bios-db-rw
ExecStart=@prefix@/bin/bios-agent-ms
Restart=always

[Install]
WantedBy=bios.target
