# src/fty-metric-store.service.in
[Unit]
Description=42Ity persistance for metrics
Requires=malamute.service bios-db-init.service
After=malamute.service bios-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/bios/etc/default/fty
EnvironmentFile=-@prefix@/share/bios/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
EnvironmentFile=-@sysconfdir@/default/bios-db-rw
ExecStart=@prefix@/bin/fty-metric-store @sysconfdir@/@PACKAGE@/fty-metric-store.cfg
Restart=always

[Install]
WantedBy=bios.target
