# src/fty-metric-store-cleaner.service.in
[Unit]
Description=42Ity metric persistance cleaner
Requires=bios-db-init.service
After=bios-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Environment="prefix=@prefix@"
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/fty__fty-metric-store.service.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/fty__fty-metric-store.service.conf
EnvironmentFile=-/etc/default/bios-db-rw
ExecStart=@prefix@/bin/fty-metric-store-cleaner @sysconfdir@/@PACKAGE@/fty-metric-store.cfg

[Install]
WantedBy=bios.target
