# src/fty-metric-store-cleaner.service.in
[Unit]
Description=42Ity metric persistance cleaner
Requires=bios-db-init.service
After=bios-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__fty-metric-store.service.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty__fty-metric-store.service.conf
Environment="prefix=@prefix@"
EnvironmentFile=-@sysconfdir@/default/bios-db-rw
ExecStart=@prefix@/bin/fty-metric-store-cleaner @sysconfdir@/@PACKAGE@/fty-metric-store.cfg

[Install]
WantedBy=bios.target
