# src/fty-metric-store-cleaner.timer.in
[Unit]
Description=Timer to call fty-metric-store-cleaner.service
# For the timer unit, we just need it to tick. Until the EULA wizard is passed,
# it should however suffice that the timer runs and occasionally tries to start
# the service unsuccessfully once a day or so...
After=fty-db.target
PartOf=bios.target

[Timer]
### Run every midnight
OnCalendar=*-*-* 23:37:00
# Do not record last-execution times
Persistent=false

[Install]
WantedBy=bios.target
